# docker/services.d/named/run
#!/bin/sh
exec /usr/sbin/named -c /etc/bind/named.conf -g -u named
